# deployment/lca-tyre.service.template
[Unit]\nDescription=LCA_TYRE dashboard\n